# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# All program logs go to fd/1 and fd/2 so they show up in `docker logs`.
# (file logging hides everything behind /var/log/supervisor/*.{log,err} which
# is invisible from the host without docker exec.)

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1280x800x24 -nolisten tcp
autorestart=true
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:x11vnc]
command=/usr/bin/x11vnc -display :1 -forever -shared -rfbport %(ENV_X11VNC_PORT)s -passwd %(ENV_VNC_PASSWORD)s -quiet
autorestart=true
startretries=10
startsecs=2
priority=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:novnc]
command=/usr/bin/websockify --web=/usr/share/novnc/ %(ENV_NOVNC_PORT)s localhost:%(ENV_X11VNC_PORT)s
autorestart=true
startsecs=2
priority=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:install_mt5]
command=/app/scripts/install_mt5.sh
autorestart=false
startsecs=0
exitcodes=0
priority=40
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:mt5linux]
command=/app/scripts/start_mt5linux.sh
autorestart=true
startretries=20
startsecs=5
priority=50
environment=DISPLAY=":1",WINEPREFIX="/config/.wine",WINEDEBUG="-all"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:tcp_forward_8001]
command=/bin/sh -c "exec socat TCP-LISTEN:8001,reuseaddr,fork TCP:127.0.0.1:%(ENV_MT5_PORT)s"
autorestart=true
startsecs=1
priority=55
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:api]
command=python3 -m uvicorn api.main:app --host 0.0.0.0 --port %(ENV_API_PORT)s --log-level info
directory=/app
autorestart=true
startretries=10
startsecs=2
priority=60
environment=PYTHONUNBUFFERED="1"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
