# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# ---- Display stack ----

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1280x800x24 -nolisten tcp
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb.err

[program:x11vnc]
command=/usr/bin/x11vnc -display :1 -forever -shared -rfbport %(ENV_X11VNC_PORT)s -passwd %(ENV_VNC_PASSWORD)s -quiet
autorestart=true
startretries=10
startsecs=2
priority=20
stdout_logfile=/var/log/supervisor/x11vnc.log
stderr_logfile=/var/log/supervisor/x11vnc.err

[program:novnc]
command=/usr/bin/websockify --web=/usr/share/novnc/ %(ENV_NOVNC_PORT)s localhost:%(ENV_X11VNC_PORT)s
autorestart=true
startsecs=2
priority=30
stdout_logfile=/var/log/supervisor/novnc.log
stderr_logfile=/var/log/supervisor/novnc.err

# ---- MT5 install (one-shot, exits when done) ----

[program:install_mt5]
command=/app/scripts/install_mt5.sh
autorestart=false
startsecs=0
exitcodes=0
priority=40
stdout_logfile=/var/log/supervisor/install_mt5.log
stderr_logfile=/var/log/supervisor/install_mt5.err

# ---- Wine-side mt5linux server (RPyC) ----
# Bound to 0.0.0.0 inside the container so both the FastAPI shim AND any
# external Python client (via the published :8001 port) can reach it.

[program:mt5linux]
command=/app/scripts/start_mt5linux.sh
autorestart=true
startretries=20
startsecs=5
priority=50
environment=DISPLAY=":1",WINEPREFIX="/config/.wine",WINEDEBUG="-all"
stdout_logfile=/var/log/supervisor/mt5linux.log
stderr_logfile=/var/log/supervisor/mt5linux.err

# ---- FastAPI REST shim ----

[program:api]
command=python3 -m uvicorn api.main:app --host 0.0.0.0 --port %(ENV_API_PORT)s --log-level info
directory=/app
autorestart=true
startretries=10
startsecs=2
priority=60
environment=PYTHONUNBUFFERED="1"
stdout_logfile=/var/log/supervisor/api.log
stderr_logfile=/var/log/supervisor/api.err

[program:tcp_forward_8001]
# socat forwards external :8001 → internal mt5linux :%(ENV_MT5_PORT)s.
# Lets external Python clients use mt5linux directly while the shim uses localhost.
command=/bin/sh -c "exec socat TCP-LISTEN:8001,reuseaddr,fork TCP:127.0.0.1:%(ENV_MT5_PORT)s"
autorestart=true
startsecs=1
priority=55
stdout_logfile=/var/log/supervisor/tcp_forward.log
stderr_logfile=/var/log/supervisor/tcp_forward.err
